# deploy/celeryd.ini
; ==================================
;  celery worker supervisor example
; ==================================

[program:celery]
; Set full path to celery program if using virtualenv
process_name = %(program_name)s_%(process_num)02d
command=/srv/rndgui/venv/bin/celery worker -E -A rndgui --loglevel=info -n worker.%%h

; Alternatively,
;command=celery --app=your_app.celery:app worker --loglevel=INFO -n worker.%%h
; Or run a script
;command=celery.sh

directory=/srv/rndgui
user=uwsgi
group=uwsgi
numprocs=4
stdout_logfile = /srv/rndgui/logs/supervisor-out.log
stderr_logfile = /srv/rndgui/logs/supervisor-err.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000
